# .config/systemd/user/auto-display.service
[Unit]
Description=Auto switch display eDP and HDMI

[Service]
Type=oneshot
ExecStart=%h/.local/bin/autodisplay_ubuntu.sh
