# .config/systemd/user/auto-display.timer
[Unit]
Description=Run auto-display every 5 seconds

[Timer]
OnBootSec=5s
OnUnitActiveSec=5s
AccuracySec=1s

[Install]
WantedBy=timers.target
