# .config/systemd/user/random-wallpaper.service
[Unit]
Description=Randomly changes wallpapers

[Service]
Type=oneshot
ExecStart=%h/.local/bin/random_wallpaper.sh
