# .config/systemd/user/random-wallpaper.service
[Unit]
Description=Randomly changes wallpapers

[Service]
Type=oneshot
ExecStart=%h/.config/i3/random_wallpaper.sh
